# server/plasma-krdp_server.service.in
[Unit]
Description=KRDP Server
After=plasma-kwin_wayland.service plasma-kcminit.service

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/krdpserver &
ExecStop=killall krdpserver
Restart=on-abnormal

[Install]
WantedBy=plasma-core.target
