# server/plasma-krdp_server.service.in
[Unit]
Description=KRDP Server
After=plasma-kwin_wayland.service plasma-kcminit.service

[Service]
Type=exec
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/krdpserver
Restart=on-abnormal

[Install]
WantedBy=plasma-core.target
