# server/plasma-krdp_server.service.in
[Unit]
Description=KRDP Server

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/krdpserver
Type=simple
